# dist/systemd/system/greenboot-post.service
[Unit]
Description=Greenboot TODO
After=boot-complete.target multi-user.target
Requires=boot-complete.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/greenboot success

[Install]
WantedBy=multi-user.target