# dist/systemd/system/greenboot-trigger.service
[Unit]
Description=Greenboot - TODO 2
DefaultDependencies=no
Conflicts=shutdown.target
Before=shutdown.target

Wants=local-fs.target
After=local-fs.target

Before=multi-user.target systemd-update-done.service
ConditionNeedsUpdate=|/etc
ConditionNeedsUpdate=|/var

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/libexec/greenboot/greenboot stamp
Restart=no

[Install]
WantedBy=multi-user.target
