# dist/systemd/system/greenboot.service
[Unit]
Description=Greenboot - TODO
After=multi-user.target
Before=boot-complete.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/greenboot/greenboot check
Restart=no

[Install]
RequiredBy=boot-complete.target
